# crontabs/root
55 20 * * * /usr/local/bin/run-geoipupdate
